# tats.service
[Unit]
Description=Hand_Tats bot

[Service]
ExecStart=/opt/handtats/run.sh
Restart=on-abort

[Install]
WantedBy=multi-user.target
